# crontest.cron
* * * * * python /Users/dongshu/PycharmProjects/spider/notify.py
