# crontest.cron
* * * * * python /Users/dongshu/PycharmProjects/spider/notify.py
* * * * * /Library/Frameworks/Python.framework/Versions/3.6/bin/python3 /Users/dongshu/PycharmProjects/spider/blog.py
