# crontest.cron
* * * * * python /home/ubuntu/zb/btc_trade_zb/notify.py
* * * * * python /home/ubuntu/zb/btc_trade_zb/blog.py
